# supertonic_tts/rootfs/etc/services.d/supertonic/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Supertonic2 TTS Wyoming Server
# ==============================================================================

bashio::log.info "Starting Supertonic2 TTS Wyoming Server..."

# Get configuration
DEFAULT_LANG=$(bashio::config 'default_language')
DEFAULT_VOICE=$(bashio::config 'default_voice')
SPEED=$(bashio::config 'speed')
VOLUME=$(bashio::config 'volume_boost')
QUALITY=$(bashio::config 'quality')
DEBUG=$(bashio::config 'debug_logging')

bashio::log.info "Configuration:"
bashio::log.info "  Default Language: ${DEFAULT_LANG}"
bashio::log.info "  Default Voice: ${DEFAULT_VOICE}"
bashio::log.info "  Speed: ${SPEED}x"
bashio::log.info "  Volume Boost: ${VOLUME}x"
bashio::log.info "  Quality: ${QUALITY}/10"
bashio::log.info "  Wyoming Protocol: tcp://0.0.0.0:10300"

# Build command arguments
ARGS=(
    --uri "tcp://0.0.0.0:10300"
    --data-dir "/data"
    --models-dir "/opt/supertonic/models"
)

# Add debug flag if enabled
if bashio::var.true "${DEBUG}"; then
    bashio::log.info "  Debug Logging: enabled"
    ARGS+=(--debug)
fi

# Run the Wyoming TTS server
exec python3 /usr/bin/supertonic_server.py "${ARGS[@]}"
